# applications/luci-app-dae/root/etc/uci-defaults/90_dae
#!/bin/sh

touch /etc/config/dae
uci set dae.config.enabled='0'
uci set dae.config.config_file='/etc/dae/config.dae'
uci commit

return 0
